# deploy/asap-kr-sync-prod.service
[Unit]
Description=ASAP KR-Sync Application (Production)
After=docker.service postgresql.service
Requires=docker.service

[Service]
Type=simple
Restart=always
RestartSec=10

ExecStartPre=-/usr/bin/docker stop asap-kr-sync-prod
ExecStartPre=-/usr/bin/docker rm asap-kr-sync-prod

ExecStart=/usr/bin/docker run --rm --name asap-kr-sync-prod \
  --add-host=host.docker.internal:host-gateway \
  --env-file /opt/asap-kr-sync-prod/.env \
  -e AWS_REGION=ca-central-1 \
  -e AUTH0_SECRET_ID=asap-kr-sync/prod/auth0 \
  -v /opt/asap-kr-sync-prod/credentials:/app/credentials:ro \
  -v /opt/asap-kr-sync-prod/backend/data:/app/src/backend/data:ro \
  -v /opt/asap-kr-sync-prod/src/frontend/public/demo-files:/app/src/frontend/public/demo-files:ro \
  -v /opt/asap-kr-sync-prod/logs:/app/logs \
  -p 127.0.0.1:3001:3000 \
  <ECR_REGISTRY>/asap-kr-sync:prod-latest

ExecStop=/usr/bin/docker stop asap-kr-sync-prod

[Install]
WantedBy=multi-user.target
